# edge_systemd/xorb-edge-resilience.service
[Unit]
Description=Xorb Edge Resilience Manager
Documentation=https://docs.xorb.ai/edge-resilience
After=network.target xorb-edge.service
Wants=network-online.target
PartOf=xorb-edge.service

[Service]
Type=notify
User=pi
Group=pi
WorkingDirectory=/opt/xorb

# Resilience manager
ExecStart=/opt/xorb/scripts/resilience_manager.py
ExecStop=/bin/kill -TERM $MAINPID
ExecReload=/bin/kill -USR1 $MAINPID

# Health monitoring
ExecStartPre=/opt/xorb/scripts/pre_start_health_check.sh
ExecStartPost=/opt/xorb/scripts/post_start_validation.sh

# Restart configuration for resilience
Restart=always
RestartSec=15
StartLimitInterval=600
StartLimitBurst=3

# Failure handling
OnFailure=xorb-edge-recovery.service

# Resource constraints
MemoryMax=1G
CPUQuota=100%
IOWeight=100

# Security
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/xorb /var/log
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictNamespaces=yes
RestrictSUIDSGID=yes
LockPersonality=yes

# Environment
Environment=PYTHONPATH=/opt/xorb
Environment=XORB_ENV=edge
Environment=XORB_RESILIENCE_MODE=enabled
Environment=XORB_LOG_LEVEL=INFO

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xorb-edge-resilience

# Watchdog for this service
WatchdogSec=120
NotifyAccess=main

[Install]
WantedBy=multi-user.target