# edge_systemd/xorb-edge-watchdog.service
[Unit]
Description=Xorb Edge Watchdog Service
Documentation=https://docs.xorb.ai/edge-deployment
After=network.target
Wants=network-online.target
Before=xorb-edge.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/xorb

# Watchdog executable
ExecStart=/opt/xorb/scripts/edge_watchdog.py
ExecStop=/bin/kill -TERM $MAINPID
ExecReload=/bin/kill -HUP $MAINPID

# Restart configuration
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Resource limits
MemoryMax=512M
CPUQuota=50%
TasksMax=50

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/xorb
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictNamespaces=yes

# Environment
Environment=PYTHONPATH=/opt/xorb
Environment=XORB_ENV=edge
Environment=XORB_NODE_TYPE=pi5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xorb-edge-watchdog

# Watchdog
WatchdogSec=60
NotifyAccess=main

[Install]
WantedBy=multi-user.target