# edge_systemd/xorb-edge-recovery.service
[Unit]
Description=Xorb Edge Recovery Service
Documentation=https://docs.xorb.ai/edge-recovery
After=network.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/opt/xorb

# Recovery script
ExecStart=/opt/xorb/scripts/emergency_recovery.sh

# Don't restart automatically - this is a recovery service
Restart=no
RemainAfterExit=no

# Extended timeout for recovery operations
TimeoutStartSec=300
TimeoutStopSec=60

# Resource limits (recovery may need more resources)
MemoryMax=2G
CPUQuota=200%

# Security (less restrictive for recovery operations)
NoNewPrivileges=no
PrivateTmp=yes
ProtectHome=yes
ReadWritePaths=/opt/xorb /var/log /etc/systemd/system

# Environment
Environment=XORB_ENV=edge
Environment=XORB_RECOVERY_MODE=emergency
Environment=XORB_LOG_LEVEL=DEBUG

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xorb-edge-recovery

[Install]
WantedBy=multi-user.target